# res/DEBIAN/postinst
#!/bin/bash

set -e

if [ "$1" = configure ]; then
	ln -f -s /usr/share/rustadmin/rustadmin /usr/bin/rustadmin

	if [ -d /run/systemd/system ] && command -v systemctl >/dev/null 2>&1; then
		rm -f /etc/systemd/system/rustadmin.service /usr/lib/systemd/user/rustadmin.service
		install -Dm644 /usr/share/rustadmin/files/systemd/rustadmin.service \
			/usr/lib/systemd/system/rustadmin.service
		systemctl daemon-reload
		# Clear a persisted stop-service flag before starting. Without this, a
		# service stopped from the UI can disable itself again during an upgrade.
		/usr/bin/rustadmin --install-service
		systemctl enable --now rustadmin.service
	fi
fi
